# etc/supervisor.d/postgresql.conf.ini
[program:postgresql]
user=postgres
command=bash -c "${PGBIN}/postgres"
